# templates/lego.timer
[Unit]
Description=LetsEncrypt renew 

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
